# supervisor.conf
[program:llm-proxy]
command=/home/w3c_offical/.local/bin/mitmdump --mode regular@18080 --ssl-insecure -s addons.py --set block_global=false --set confdir=/home/w3c_offical/.mitmproxy
directory=/home/w3c_offical/projects/ai-workers/llm-proxy
autostart=true
autorestart=true
stdout_logfile=/home/w3c_offical/projects/ai-workers/llm-proxy/logs/proxy.out.log
stderr_logfile=/home/w3c_offical/projects/ai-workers/llm-proxy/logs/proxy.err.log
user=w3c_offical
